# root/etc/systemd/system/getty@tty8.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --autologin mythfrontend tty8 linux
