# root/etc/systemd/system/getty@tty8.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/openvt -fwc 8 -- /bin/su - mythfrontend -c /usr/bin/startx >& /home/mythfrontend/mfout
